# deploy/supervisor/city-parser-gun-az.conf
; city-parser.gun.az — bu layihə üçün ayrıca Horizon (parser.gun.az və s. ilə qarışmır).
; Copy: sudo cp deploy/supervisor/city-parser-gun-az.conf /etc/supervisor/conf.d/
; Then: sudo supervisorctl reread && sudo supervisorctl update
;
; ChromeDriver: ümumi chromedriver.conf ola bilər; bu layihənin .env CHROME_DRIVER_URL uyğun portu göstərməlidir.
;
; Scheduler: cron (hər layihə üçün ayrıca sətir). www-data:
;   * * * * * cd /var/www/city-parser.gun.az && php artisan schedule:run >> /dev/null 2>&1
; schedule:work Supervisor-da bu layihə üçün ayrıca saxlamayın əgər cron istifadə edirsinizsə.

[program:city-parser-horizon]
process_name=%(program_name)s
command=php /var/www/city-parser.gun.az/artisan horizon
directory=/var/www/city-parser.gun.az
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/city-parser.gun.az/storage/logs/horizon.log
stopwaitsecs=3600
